# packaging/systemd/user/myapp-input-daemon.service
[Unit]
Description=MyApp Input Daemon

[Service]
ExecStart=/usr/bin/myapp-daemon
Restart=on-failure

[Install]
WantedBy=default.target

